# root/etc/services.d/webhook/run
#!/usr/bin/with-contenv sh

cd /etc/webhook || exit 1;
HOOKPARAMS=${HOOKPARAMS:- -verbose}

s6-setuidgid alpine webhook -hotreload -hooks=/etc/webhook/hooks.json ${HOOKPARAMS}
